# services/bluetooth-agent.service
[Unit]
Description=Bluetooth Pairing Agent
Requires=bluetooth.service
After=bluetooth.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/bluetooth-receiver
ExecStart=/opt/bluetooth-receiver/venv/bin/python /opt/bluetooth-receiver/bt_agent.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
