# services/usb-player.service
[Unit]
Description=Local Music Player
After=network.target sound.target
Wants=sound.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/bluetooth-receiver
ExecStart=/opt/bluetooth-receiver/venv/bin/python /opt/bluetooth-receiver/usb_player.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
